# build/linux/debian/etc/init.d/bitmeter
#!/bin/sh

### BEGIN INIT INFO
# Provides:          BitMeter_Capture_Daemon
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts instance of BitMeter Capture Daemon
# Description:       starts instance of BitMeter Capture Daemon using start-stop-daemon
### END INIT INFO

# http://www.debian.org/doc/debian-policy/ch-opersys.html
DAEMON=/usr/bin/bmcapture
PROC=bmcapture

test -f $DAEMON || exit 0

doStart(){
    echo -n Starting BitMeter Capture daemon: bmcapture
    start-stop-daemon --start --quiet --oknodo --background --name $PROC --exec $DAEMON
    echo .
}
doStop(){
    echo -n Stopping BitMeter Capture daemon: bmcapture
    start-stop-daemon --stop --quiet --signal 2 --oknodo --exec $DAEMON
    echo .    
}
case "$1" in
  start)
    doStart
    ;;
  stop)
    doStop
    ;;
  restart)
    doStop
    doStart
    ;;
  force-reload)
    doStop
    doStart
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}" >&2
    exit 1
    ;;
esac

exit 0

